# godaddy_ddns_update.timer
[Unit]
Description=Run Godaddy DDNS IP update every 10 minutes

[Timer]
OnCalendar=*-*-* *:*:10
Persistent=true

[Install]
WantedBy=timers.target
